# backend/deploy/systemd/mrsm-prometheus-token-refresh.timer
[Unit]
Description=Run MRSM Prometheus token refresh every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=mrsm-prometheus-token-refresh.service

[Install]
WantedBy=timers.target
